# jupyter.service
[Unit]
Description=jupyter notebook service
After=network.target

[Service]
ExecStart=/bin/bash /home/pi/Desktop/PinoBot/runJupyter.sh
WorkingDirectory=/home/pi/Desktop/PinoBot
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target